# roles/buildslave/files.fedora18-64/etc/rc.d/rc.local
#!/bin/sh -e

su - buildbot -s /bin/sh \
  -c "cd /srv/buildbot/slaves/unittests-fedora18-64 && buildslave start -q"

exit 0
