# roles/buildslave/files.fedora17-64/etc/rc.d/rc.local
#!/bin/sh -e

su - buildbot -s /bin/sh \
  -c "cd /srv/buildbot/slaves/unittests-fedora17-64 && buildslave start -q"

exit 0
