# roles/buildslave/files.ubuntu1304-64/etc/rc.local
#!/bin/sh -e

su - buildbot -s /bin/sh \
  -c "cd /srv/buildbot/slaves/unittests-ubuntu1304-64 && buildslave start -q"

exit 0
